# autostart/keyboard_layout.sh.desktop
[Desktop Entry]
Exec=~/dotfiles/keyboard_layout.sh
Icon=dialog-scripts
Name=keyboard_layout.sh
Path=
Type=Application
X-KDE-AutostartScript=true
